# samples/web_app_basic/config/flask_project.conf
[program:flask_project]
environment = PYTHONPATH = /home/www/env/lib/python2.7/site-packages
command = /usr/bin/gunicorn app:app -b localhost:8000
directory = /home/www/flask_project
user = vagrant
